# deployment/supervisor/eservice-worker.conf
[program:eservice-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/eservice-app/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --timeout=60
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/eservice-app/storage/logs/worker.log
stopwaitsecs=3600

; Konfigurasi ini untuk Laravel Queue Worker
; 
; CARA INSTALL:
; 1. Copy file ini ke server production: /etc/supervisor/conf.d/eservice-worker.conf
; 2. Sesuaikan path dan user dengan server Anda
; 3. Jalankan: sudo supervisorctl reread
; 4. Jalankan: sudo supervisorctl update
; 5. Jalankan: sudo supervisorctl start eservice-worker:*
;
; CATATAN PENTING:
; - Sesuaikan path: /var/www/eservice-app dengan path aplikasi Anda
; - Sesuaikan user: www-data dengan user web server Anda (bisa juga nginx, apache, dll)
; - Sesuaikan numprocs: 2 adalah jumlah worker yang berjalan paralel
; - Log akan tersimpan di: /var/www/eservice-app/storage/logs/worker.log
;
; SETELAH DEPLOYMENT CODE BARU:
; sudo supervisorctl restart eservice-worker:*
